# roles/mymon/templates/mymon.cron.j2
* * * * * root /bin/bash {{ mymon_dir }}/mymon_control.sh
